# make/cacerts/files/root/etc/profile.d/cacerts.sh
# setting CA-Bundle file location via SSL_CERT_FILE - used by OpenSSL-lib
export SSL_CERT_FILE=/etc/ssl/certs/cacert.pem
