# make/ca-bundle/files/root/etc/profile.d/ca-bundle.sh
export SSL_CERT_FILE=/mod/etc/ssl/certs/ca-bundle.crt

