# packages/mini_fo-0.2b/root/etc/inittab
::restart:/sbin/init

# Start mini_fo before rc.S
::sysinit:/etc/init.d/rc.mini_fo
::sysinit:/etc/init.d/rc.S

# Start an "askfirst" shell on the console (whatever that may be)
ttyS0::askfirst:-/bin/sh

# Stuff to do before rebooting
::shutdown:/bin/sh -c /var/post_install
